# opt/opt-post-clone/my-rc-local.service
[Unit]
Description=My rc.local

[Service]
ExecStart=/usr/local/bin/rc-local.sh
Restart=always

[Install]
WantedBy=multi-user.target
