# opt/opt-post-clone/my-rc-local.service
[Unit]
Description=My rc.local

[Service]
ExecStart=/usr/local/bin/on-start.sh

[Install]
WantedBy=multi-user.target
