# opt/systemd/my-rc-local.service
[Unit]
Description=My rc.local

[Service]
ExecStart=/usr/local/bin/on-start.sh
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
