# opt/my-rc-local.service
[Unit]
Description=My rc.local

[Service]
ExecStart=/usr/local/bin/on-start.sh
Restart=always

[Install]
WantedBy=multi-user.target
